# service/OpenRC/init.d/shortcutd
#!/sbin/openrc-run

name="shortcut"
description="Shortcut file caching web-proxy"
command="/opt/bin/shortcutd"
pidfile="/run/shortcutd.pid"
log_file="/var/log/shortcutd.log"

depend() {
    need net
    use logger
}

start_pre() {
    # Check for the mandatory environment variable
    if [ -z "${X_PROXY_CACHE_PATH}" ]; then
        eerror "Environment variable X_PROXY_CACHE_PATH is not set"
        return 1
    fi
    checkpath --directory --owner ${SHORTCUT_COMMAND_USER}:${SHORTCUT_COMMAND_GROUP} ${X_PROXY_CACHE_PATH}
    checkpath --file --owner ${SHORTCUT_COMMAND_USER}:${SHORTCUT_COMMAND_GROUP} ${log_file}
}

start() {
    ebegin "Starting ${name}"

    export X_PROXY_CACHE_PATH="${X_PROXY_CACHE_PATH}"

    start-stop-daemon --start --exec ${command} \
                      --user ${SHORTCUT_COMMAND_USER} \
                      --group ${SHORTCUT_COMMAND_GROUP} \
                      --pidfile ${pidfile} \
                      --make-pidfile \
                      --background \
                      --stdout ${log_file} \
                      --stderr ${log_file} \
                      ${command_args}
    eend $?
}

stop() {
    ebegin "Stopping ${name}"
    start-stop-daemon --stop --pidfile ${pidfile}
    eend $?
}